# roles/kea_dhcp/templates/webhook.service.j2
[Unit]
Description=Webhook Service

[Service]
WorkingDirectory={{ script_dir }}
ExecStart=/usr/bin/webhook -hooks hooks.json

[Install]
WantedBy=multi-user.target
